# deploy/key.service
[Unit]
Description=Keys API Key Manager (FastAPI)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=user1
Group=user1
WorkingDirectory=/home/user1/keys
EnvironmentFile=/home/user1/keys/.env

# 首次部署（必须用 user1 执行）: cd /home/user1/keys && uv sync
ExecStart=/home/user1/keys/.venv/bin/python -m uvicorn app.main:app --host 0.0.0.0 --port 8888 --workers 1

Restart=on-failure
RestartSec=5
TimeoutStopSec=30

ReadWritePaths=/home/user1/keys

[Install]
WantedBy=multi-user.target
