# deploy/key.service
[Unit]
Description=Keys API Key Manager (FastAPI)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/root/keys
EnvironmentFile=/root/keys/.env

# 首次部署: cd /root/keys && uv sync
ExecStart=/root/keys/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8888 --workers 1

Restart=on-failure
RestartSec=5
TimeoutStopSec=30

ReadWritePaths=/root/keys

[Install]
WantedBy=multi-user.target
